# dot_config/systemd/user/backup.timer
[Unit]
Description=Backup my homedir 15 minutes after logging in

[Timer]
OnBootSec=15min

[Install]
WantedBy=timers.target
