# dot_config/systemd/user/discord-delete.timer
[Unit]
Description=Timer for discord-delete.service

[Timer]
OnCalendar=weekly

[Install]
WantedBy=timers.target
