# dot_config/systemd/user/free-games-claimer.timer
[Unit]
Description=Timer for free-games-claimer.service

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
